# configs/update-mpd.service
[Unit]
Description = Update MPD

[Service]
ExecStart = /usr/local/bin/update-mpd

[Install]
WantedBy = multi-user.target
